# etc/supervisor/conf.d/server.conf
[program:server]
command = node /root/server.js
directory = /root/
autostart = true
autorestart = true
priority = 1
stdout_logfile = /root/server.log